# websocket-server/websocket-transfer.service
[Unit]
Description=WebSocket Transfer Server
Documentation=https://github.com/your-repo/websocket-transfer-server
After=network.target
Wants=network.target

[Service]
Type=simple
User=websocket
Group=websocket
WorkingDirectory=/opt/websocket-transfer-server/websocket-server
ExecStart=/opt/websocket-transfer-server/websocket-server/venv/bin/python server.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=websocket-transfer

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/websocket-transfer-server/websocket-server
ReadWritePaths=/var/log/websocket-transfer

# リソース制限
LimitNOFILE=65536
LimitNPROC=4096

# 環境変数
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/websocket-transfer-server/websocket-server

[Install]
WantedBy=multi-user.target
